# lcm/bios_configs/load-generator.conf
[program: load-generator]
directory = /opt/bios
command = /opt/bios/apps/load-generator/start-load-generator.sh